# macos/clamav-sync/net.dpl.clamav-sync-cvd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>net.dpl.clamav-sync-cvd</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-x</string>
    <string>/usr/local/sbin/clamav_sync_cvd.sh</string>
  </array>

  <key>RunAtLoad</key><true/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>3</integer>  <!-- Tue -->
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>KeepAlive</key>
  <dict><key>NetworkState</key><true/></dict>

  <key>StandardOutPath</key><string>/var/log/net.dpl.clamav-sync-cvd.log</string>
  <key>StandardErrorPath</key><string>/var/log/net.dpl.clamav-sync-cvd.log</string>

  <key>WorkingDirectory</key><string>/</string>
  <key>ProcessType</key><string>Background</string>
</dict>
</plist>
